# 90-autorun.rules
ACTION=="add",KERNEL=="sda1",RUN+="/home/pi/automount/mounting_script.sh %k"
