# 90-autorun.rules
ACTION=="add",KERNEL=="sd*1",RUN+="/usr/share/automount/mounting_script.sh %k"
