# 90-autorun.rules
ACTION=="add",KERNEL=="sda1",RUN+="/usr/share/automount/mounting_script.sh %k"
